# packaging/linux/debian/prerm
#!/bin/bash
# Pre-remocion de Kipu
# Se ejecuta automaticamente antes de desinstalar el paquete .deb

set -e

# Remover symlinks
rm -f /usr/local/bin/kipu 2>/dev/null || true
rm -f /usr/local/bin/kipu-cliente 2>/dev/null || true
rm -f /usr/local/bin/kipu-servidor 2>/dev/null || true

# Detener servidor si esta corriendo
if command -v systemctl &> /dev/null; then
    systemctl stop kipu-servidor 2>/dev/null || true
    systemctl disable kipu-servidor 2>/dev/null || true
fi

# Intentar matar procesos de Kipu que sigan corriendo
pkill -f "KipuServidor" 2>/dev/null || true
pkill -f "kipu-servidor" 2>/dev/null || true

echo "Kipu desinstalado."
